# ROOT/btrfs/etc/systemd/system/btrfs-snapshot.timer
[Unit]
Description=Daily Btrfs snapshots

[Timer]
OnCalendar=06:00
AccuracySec=5m
Persistent=true

[Install]
WantedBy=timers.target
